# src/presistensy/prince.service
[Unit]
Description=prince service
[Service]
ExecStart=/home/yael/snap/code/current/pres/tmp