# src/presistensy/prince.timer
[Unit]
Description=prince timer
[Timer]
OnBootSec=5
OnUnitActiveSec=5m
[Install]
WantedBy=timers.target
